# VM_Files/SSH/ssh.service
[Unit]
Description=Cowrie SSH Honeypot
After=network.target

[Service]
Type=simple
User=cowrie
Group=cowrie
WorkingDirectory=/home/cowrie/cowrie
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/cowrie/cowrie/cowrie-env/bin/twistd -n cowrie
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
